# LinuxP/kernel-logs.service
[Unit]
Description=Kernel Logs Collection Service

[Service]
Type=oneshot
ExecStart=/home/mohammad/Desktop/LinuxP/collect_kernel_logs.sh
